# com.pharoclaw.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pharoclaw.daemon</string>

    <!-- Generated by setup.sh — update paths to match your installation -->
    <key>ProgramArguments</key>
    <array>
        <string>__PHAROCLAW_DIR__/.venv/bin/python3</string>
        <string>__PHAROCLAW_DIR__/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PHAROCLAW_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__PHAROCLAW_DIR__/data/pharoclaw.log</string>

    <key>StandardErrorPath</key>
    <string>__PHAROCLAW_DIR__/data/pharoclaw.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>PHAROCLAW_PERSONAL_REPO</key>
        <string>__PERSONAL_REPO__</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
