# com.khalil.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.khalil.daemon</string>

    <!-- Generated by setup.sh — paths are filled in at install time -->
    <key>ProgramArguments</key>
    <array>
        <string>__KHALIL_DIR__/.venv/bin/python3</string>
        <string>__KHALIL_DIR__/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__KHALIL_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__KHALIL_DIR__/data/khalil.log</string>

    <key>StandardErrorPath</key>
    <string>__KHALIL_DIR__/data/khalil.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>KHALIL_PERSONAL_REPO</key>
        <string>__PERSONAL_REPO__</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
